# Nutriet/nutriet.service
# ══════════════════════════════════════════════════════════════════════════════
#  NutriET — Servicio systemd (para VPS / servidor Linux dedicado)
#  Ubicación: /etc/systemd/system/nutriet.service
#  Activar:
#    sudo systemctl daemon-reload
#    sudo systemctl enable nutriet
#    sudo systemctl start nutriet
#    sudo systemctl status nutriet
# ══════════════════════════════════════════════════════════════════════════════

[Unit]
Description=NutriET Django Application (Gunicorn)
After=network.target postgresql.service redis.service
Wants=postgresql.service

[Service]
Type=notify
User=www-data
Group=www-data

# ── Rutas (ajustar según tu servidor) ────────────────────────────────────────
WorkingDirectory=/var/www/nutriet
ExecStart=/var/www/nutriet/venv/bin/gunicorn \
    -c /var/www/nutriet/gunicorn.conf.py \
    Nutriet.wsgi:application

# ── Variables de entorno ──────────────────────────────────────────────────────
EnvironmentFile=/var/www/nutriet/.env
Environment="DJANGO_SETTINGS_MODULE=Nutriet.setting.prod"

# ── Reinicio automático ───────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5s

# ── Límites de recursos ───────────────────────────────────────────────────────
LimitNOFILE=65536
TimeoutStartSec=60
TimeoutStopSec=30

# ── Logs (ver con: journalctl -u nutriet -f) ──────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nutriet

[Install]
WantedBy=multi-user.target
